# debian/postinst
#!/bin/sh
set -e

replace_or_append() {
  key="$1"
  value="$2"
  file="$3"
  if [ ! -f "$file" ]; then
    printf "%s=%s\n" "$key" "$value" > "$file"
    return 0
  fi
  if grep -q "^${key}=" "$file"; then
    sed -i "s|^${key}=.*|${key}=${value}|" "$file"
  else
    printf "%s=%s\n" "$key" "$value" >> "$file"
  fi
}

if [ "$1" = "configure" ] || [ "$1" = "reconfigure" ]; then
  mkdir -p /opt/fabricator-agent /etc/fabricator-agent

  if [ ! -f /etc/fabricator-agent/config.toml ]; then
    cat > /etc/fabricator-agent/config.toml <<'EOF'
# Optional local agent config
EOF
    chmod 640 /etc/fabricator-agent/config.toml || true
  fi

  # Migrate old default that forced local test mode and blocked pairing.
  # If AGENT_ALLOW_TEST_MODE=1 is explicitly set, keep AGENT_TEST_MODE as-is.
  if [ -f /etc/default/fabricator-agent ] && grep -q '^AGENT_ALLOW_TEST_MODE=1' /etc/default/fabricator-agent; then
    true
  else
    replace_or_append AGENT_TEST_MODE 0 /etc/default/fabricator-agent
  fi

  PKG_VERSION="$(dpkg-query -W -f='${Version}' fabricator-agent 2>/dev/null || true)"
  if [ -n "${PKG_VERSION}" ]; then
    replace_or_append FABRICATOR_AGENT_VERSION "${PKG_VERSION}" /etc/default/fabricator-agent
  fi

  python3 -m venv /opt/fabricator-agent/.venv
  /opt/fabricator-agent/.venv/bin/pip install --upgrade pip >/dev/null
  /opt/fabricator-agent/.venv/bin/pip install -r /opt/fabricator-agent/requirements.txt >/dev/null

  systemctl daemon-reload || true
  systemctl enable --now fabricator-agent.service || true
fi

exit 0
